# raspi/etc__init.d__minibot-battery-check
#! /bin/sh
### BEGIN INIT INFO
# Provides:          minibot-battery-check
# Required-Start:
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Starts & Stops the minibot battery stataus
# Description:       Observing and displaying the battery status on the OLED.
### END INIT INFO

#Switch case fuer den ersten Parameter
case "$1" in
    start)
	#Aktion wenn start uebergeben wird
	# touch /home/pi/catkin_ws/minibot-battery-check_START
        echo "Starting minibot power monitor"
        sudo python /home/pi/minibot/test/power_monitor.py &
        ;;

    stop)
	#Aktion wenn stop uebergeben wird
        echo "Stopping minibot-battery-check"
	# touch /home/pi/catkin_ws/minibot-battery-check_STOP
        ;;

    restart)
	#Aktion wenn restart uebergeben wird
        echo "Restarting minibot-battery-check"
	# touch /home/pi/catkin_ws/minibot-battery-check_RESTART
        ;;

    *)
	#Standard Aktion wenn start|stop|restart nicht passen
	echo "(start|stop|restart)"
	# touch /home/pi/catkin_ws/minibot-battery-check_REST
	;;
esac

exit 0
